# pi/drone-bridge.service
[Unit]
Description=Drone Communication Bridge
Documentation=file:///home/pi/drones_over_cellular/pi/README.md
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/drones_over_cellular/pi
ExecStart=/usr/bin/python3 /home/pi/drones_over_cellular/pi/communication_bridge.py --verbose
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=drone-bridge

# Environment variables
Environment=PYTHONPATH=/home/pi/drones_over_cellular
Environment=PYTHONUNBUFFERED=1

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=80%

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp /var/log /home/pi/drones_over_cellular

# Network settings
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=10.0.0.0/8
IPAddressAllow=192.168.0.0/16
IPAddressAllow=172.16.0.0/12

[Install]
WantedBy=multi-user.target
